# redwing-ngrok.service
[Unit]
Description=RedWing DB Automation Ngrok Tunnel
After=redwing-backend.service
Requires=redwing-backend.service

[Service]
User=jester
Group=jester
# We need to load NGROK_DOMAIN from the .env file. We can just use the config if needed or load env.
# A simple way to wrap with bash and grep the env file.
EnvironmentFile=/home/jester/Auto-DB-Updater/.env
ExecStart=/usr/local/bin/ngrok http --domain=${NGROK_DOMAIN} 8000
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
